# deploy/com.onlineworker.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.onlineworker</string>

    <key>ProgramArguments</key>
    <array>
        <string>/path/to/python3</string>
        <string>/path/to/onlineWorker/main.py</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/path/to/onlineWorker</string>

    <!-- Update these paths before loading the launchd plist. -->
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/path/to/python-bin:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
        <key>HOME</key>
        <string>/Users/your-user</string>
        <key>LANG</key>
        <string>en_US.UTF-8</string>
    </dict>

    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
    </dict>

    <key>ThrottleInterval</key>
    <integer>10</integer>

    <key>StandardOutPath</key>
    <string>/tmp/onlineworker-launchd.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/onlineworker-launchd.log</string>

    <key>RunAtLoad</key>
    <true/>

    <key>ProcessType</key>
    <string>Interactive</string>
</dict>
</plist>
